# rctc.cron
# /etc/cron.d/rctc: crontab fragment for rctc

0 1     * * *     root   /usr/bin/cronic /usr/sbin/gzip-all-tc-logs.sh
0 2     * * *     root   /usr/bin/cronic /usr/sbin/archive-all-tc-logs.sh

